# services/runner/launchd/com.searchx.autopull.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  searchX 定时自动拉（仅在常驻不关机的 Mac mini 上加载）。
  作用：每 120 秒跑一次 .claude/hooks/autopull.sh → git-sync.sh pull，把另一台/runner/网页
        push 的进度自动同步下来，无需在 Mac mini 开 Claude 窗口或手动 pull。
  安装（在 Mac mini 上）：
      cp services/runner/launchd/com.searchx.autopull.plist ~/Library/LaunchAgents/
      launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.searchx.autopull.plist
      launchctl kickstart -k gui/$(id -u)/com.searchx.autopull   # 立刻跑一次验证
  卸载：launchctl bootout gui/$(id -u)/com.searchx.autopull
  StartInterval=120 即每 2 分钟。改 60=1 分钟、300=5 分钟。
  RunAtLoad=true：装载即拉一次（拉取免费且安全，与花额度的 runner 不同）。
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.searchx.autopull</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/yangqiuyuan/Coding/searchX/.claude/hooks/autopull.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/yangqiuyuan/Coding/searchX</string>

  <key>StartInterval</key>
  <integer>120</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>ProcessType</key>
  <string>Background</string>

  <key>LowPriorityIO</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/yangqiuyuan/Library/Logs/searchx-runner/autopull.launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/yangqiuyuan/Library/Logs/searchx-runner/autopull.launchd.err.log</string>
</dict>
</plist>
